# trackpoint_parameters.service
[Unit]
Description=Set TrackPoint attributes

[Service]
ExecStart=/home/sajay/dotfiles/trackpoint_configuration.sh
